# deploy/systemd/gae-tunnel.service
[Unit]
Description=OpenAI Secure MCP Tunnel for GPT-Agent4EC2 (gae)
Documentation=https://developers.openai.com/api/docs/guides/secure-mcp-tunnels
Wants=network-online.target
After=network-online.target gae.service
Requires=gae.service
ConditionPathExists=/etc/gae-tunnel/profile.yaml
ConditionPathExists=/etc/gae-tunnel/runtime-api-key
ConditionPathExists=/etc/gae-tunnel/internal-secret

[Service]
Type=simple
User=ubuntu
Group=ubuntu
LoadCredential=runtime-api-key:/etc/gae-tunnel/runtime-api-key
LoadCredential=internal-secret:/etc/gae-tunnel/internal-secret
ExecStart=/usr/local/bin/tunnel-client run --profile-file /etc/gae-tunnel/profile.yaml --control-plane.api-key file:%d/runtime-api-key --mcp.extra-headers "X-Devspace-Internal-Key: file:%d/internal-secret" --mcp.discovery-extra-headers "X-Devspace-Internal-Key: file:%d/internal-secret"
Restart=always
RestartSec=5s
TimeoutStartSec=45s
TimeoutStopSec=20s
KillSignal=SIGTERM
UMask=0077

# Keep Minecraft and GAE ahead of the tunnel transport when resources are tight.
Nice=12
CPUWeight=10
IOWeight=10
MemoryHigh=128M
MemoryMax=256M
OOMScoreAdjust=700
TasksMax=128
LimitNOFILE=8192

# The daemon only needs outbound HTTPS and read-only access to its profile.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native
ReadOnlyPaths=/etc/gae-tunnel/profile.yaml

[Install]
WantedBy=multi-user.target
